# launchd/com.unao.cafe-crawler-popular-evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.unao.cafe-crawler-popular-evening</string>

    <key>ProgramArguments</key>
    <array>
        <!-- node 사용 이유: /bin/bash는 kTCCServiceSystemPolicyDocumentsFolder TCC 미승인 →
             launchd에서 ~/Documents/ 접근 불가. node는 TCC 승인됨. -->
        <string>/Users/yanadoo/.nvm/versions/node/v24.14.0/bin/node</string>
        <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/scripts/launchd-wrapper.mjs</string>
        <string>cafe-crawler-popular-evening</string>
        <string>/Users/yanadoo/.nvm/versions/node/v24.14.0/bin/npx</string>
        <string>tsx</string>
        <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/agents/cafe/popular-sync.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yanadoo/Documents/New_Claude_agenotmatter</string>

    <!-- 21:15 KST (popular-sync: 인기글 수집 → GHA popular-curate 21:30에 발행, 저녁 크롤 21:30 전 완료) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/logs/cafe-crawler-popular-evening.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/logs/cafe-crawler-popular-evening-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/yanadoo/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/yanadoo</string>
    </dict>
</dict>
</plist>
