# launchd/com.unaeo.magazine-late.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.unaeo.magazine-late</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/scripts/launchd-alert.sh</string>
    <string>magazine-late</string>
    <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/node_modules/.bin/tsx</string>
    <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/agents/cafe/local-magazine-runner.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yanadoo/Documents/New_Claude_agenotmatter</string>

  <!-- 14:00 KST (오후 매거진 발행, IMAGE_GENERATOR=gemini) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>14</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/logs/magazine-late.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/logs/magazine-late.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/node_modules/.bin:/Users/yanadoo/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/yanadoo</string>
    <key>SESSION_TIME</key>
    <string>late</string>
    <key>IMAGE_GENERATOR</key>
    <string>gemini</string>
    <key>NODE_OPTIONS</key>
    <string>--require /Users/yanadoo/Documents/New_Claude_agenotmatter/agents/cafe/preload-eagain-retry.cjs</string>
  </dict>
</dict>
</plist>
