# launchd/com.unaeo.magazine-late.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.unaeo.magazine-late</string>

  <key>ProgramArguments</key>
  <array>
    <!-- node 사용 이유: /bin/bash는 kTCCServiceSystemPolicyDocumentsFolder TCC 미승인 →
         launchd에서 ~/Documents/ 접근 불가 (Operation not permitted, exit 126).
         node(/Users/yanadoo/.nvm/versions/node/v24.14.0/bin/node)는 TCC 승인됨. -->
    <string>/Users/yanadoo/.nvm/versions/node/v24.14.0/bin/node</string>
    <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/scripts/launchd-wrapper.mjs</string>
    <string>magazine-late</string>
    <string>/Users/yanadoo/.nvm/versions/node/v24.14.0/bin/npx</string>
    <string>tsx</string>
    <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/agents/cafe/local-magazine-runner.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yanadoo/Documents/New_Claude_agenotmatter</string>

  <!-- 14:00 KST (오후 매거진 발행, IMAGE_GENERATOR=chatgpt) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>14</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/logs/magazine-late.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/logs/magazine-late.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/yanadoo/Documents/New_Claude_agenotmatter/node_modules/.bin:/Users/yanadoo/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/yanadoo</string>
    <key>SESSION_TIME</key>
    <string>late</string>
    <key>IMAGE_GENERATOR</key>
    <string>chatgpt</string>
    <key>NODE_OPTIONS</key>
    <string>--require /Users/yanadoo/Documents/New_Claude_agenotmatter/agents/cafe/preload-eagain-retry.cjs</string>
  </dict>
</dict>
</plist>
