# launchd/com.unaeo.magazine-evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.unaeo.magazine-evening</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-c</string>
    <string>cd /Users/yanadoo/Documents/New_Claude_agenotmatter &amp;&amp; set -a &amp;&amp; source .env.local &amp;&amp; set +a &amp;&amp; SESSION_TIME=evening IMAGE_GENERATOR=chatgpt /usr/local/bin/npx tsx agents/cafe/local-magazine-runner.ts >> /tmp/magazine-evening.log 2>&amp;1</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>21</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/magazine-evening.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/magazine-evening.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>/Users/yanadoo</string>
    <key>SESSION_TIME</key>
    <string>evening</string>
    <key>IMAGE_GENERATOR</key>
    <string>chatgpt</string>
  </dict>
</dict>
</plist>
